# weater.service
[Unit]
Description=Weather station

[Service]
Type=simple
ExecStart=/home/pi/Weather/Prognoza_sait.py &

[Install]
WantedBy=multi-user.target

